# core/systemd/aionos.target
[Unit]
Description=AIONOS Stack Target
Requires=aionos-gateway.service aionos-control.service aionos-console.service
After=network.target

[Install]
WantedBy=multi-user.target
