# deploy/headless-bundle/services/aion-gateway.service
[Unit]
Description=AION Gateway (gRPC/HTTP edge)
After=network.target aion-control.service

[Service]
User=aion
EnvironmentFile=/home/aion/.aion.env
WorkingDirectory=/opt/aionos/gateway
ExecStart=/opt/aionos/gateway/bin/aion-gateway --config /etc/aionos/config/router.policy.yaml
Restart=always
RestartSec=3
StandardOutput=append:/var/log/aionos/gateway.log
StandardError=append:/var/log/aionos/gateway.err

[Install]
WantedBy=multi-user.target
